# roles/pihole/templates/update_hosts.service.j2
[Unit]
Description=Update Hosts file blocking ads and malware
After=network.target

[Service]
ExecStart=/usr/local/bin/update_hosts
